# examples/scheduling/systemd/axil-daily-brief.service
[Unit]
Description=Axil daily brief

[Service]
Type=oneshot
# Update the path to axil if it is not on PATH for systemd user-services.
ExecStart=/bin/sh -lc 'axil brief --window 24h --brief-format markdown > /tmp/axil-daily-brief.md'
